# certbot-renew.timer
[Unit]
Description=Daily Certbot renewal timer

[Timer]
OnBootSec=10min
OnUnitActiveSec=24h
Persistent=true

[Install]
WantedBy=timers.target